# crates/agent/templates/default-forge-dhcp-server.conf
# Auto-generated by Forge!
[program: forge-dhcp-server-default]
command = /var/support/forge-dhcp/bin/forge-dhcp-server {{range .Interfaces}} --interfaces {{.}} {{ end }} --dhcp-config /var/support/forge-dhcp/conf/dhcp.yaml --host-config /var/support/forge-dhcp/conf/host.yaml
autostart = {{ .Autostart }}
autorestart = {{ .Autostart }}
startsecs = 3
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 3
stdout_logfile = /var/support/forge-dhcp/logs/forge-dhcp-server.log
redirect_stderr = true
